# lightrag/com.user.lightrag_flask.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.lightrag_flask</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>/Users/agents/.claude/start_flask_lightrag.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/agents/.claude/logs/lightrag_flask_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/agents/.claude/logs/lightrag_flask_output.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/agents/.claude</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:/Users/agents/.claude/venv/bin</string>
    </dict>
</dict>
</plist>